# gandalf-server-deb/etc/gandalf-server.upstart
limit nofile 10000 10000

kill timeout 180

start on runlevel [2345]
stop on runlevel [06]

respawn
setuid git
setgid git
exec /usr/bin/gandalf-server -config=/etc/gandalf.conf
post-stop exec sleep 5
