# archive-server-deb/debian/archive-server.upstart
limit nofile 10000 10000

kill timeout 180

start on runlevel [2345]
stop on runlevel [06]

respawn
setuid archiver
setgid archiver

script
  if [ -f /etc/default/archive-server ]; then . /etc/default/archive-server; fi
  exec /usr/bin/archive-server $ARCHIVE_SERVER_OPTS
end script

post-stop exec sleep 5
