# tsuru-server-deb/debian/tsuru-server-api.upstart
limit nofile 10000 10000

kill timeout 180

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit unlimited

setuid tsuru
setgid tsuru

script
  CONFIG_FILE="/etc/tsuru/tsuru.conf"
  if [ -f /etc/default/tsuru-server ]; then . /etc/default/tsuru-server; fi
  if [ "$TSR_API_ENABLED" = "yes" -o "$TSR_API_ENABLED" = "true" ]; then exec /usr/bin/tsurud api --config=${CONFIG_FILE}; fi
end script
